# .config/autostart/albert.desktop
[Desktop Entry]
Name=albert
Comment=albert
Exec=albert
Type=Application
Terminal=false
Hidden=false

