# .config/autostart/discord.desktop
[Desktop Entry]
Name=discord
Comment=discord start
Exec=discord
Type=Application
Terminal=false
Hidden=false

